# tripleo_ansible/roles/tripleo_systemd_wrapper/templates/service.path.j2
[Path]
PathModified={{ tripleo_systemd_wrapper_service_dir }}/{{ tripleo_systemd_wrapper_service_name }}/processes-timestamp

[Install]
WantedBy=multi-user.target
